# pkg/configurator/templates/ansible/roles/docker/registry/files/etc/systemd/system/registry.service.d/registry.conf
[Unit]
After=docker.service
Wants=docker.service
PartOf=docker.service

[Install]
WantedBy=docker.service
