# pkg/configurator/templates/ansible/roles/kubernetes/ha/files/etc/systemd/system/keepalived.service.d/keepalived.conf
[Unit]
After=bynet.service
Wants=bynet.service
PartOf=bynet.service

[Install]
WantedBy=bynet.service
